# systemd/user/needrestart-notify.service
[Unit]
Description=needrestart-notify

[Service]
Type=simple
ExecStart=%h/.dotfiles/bin/needrestart-notify

[Install]
WantedBy=default.target
